# music-daemon/supervisord.conf
[supervisord]
loglevel=debug
logfile=/var/log/supervisor/supervisord.log
nodaemon=true

[program:mpd]
command=mpd --no-daemon /etc/mpd.conf

[program:node]
directory=/app
command=node server2.js
